# prometheus-borg-exporter.timer
[Unit]
Description=Prometheus Borg Exporter Timer

[Timer]
OnBootSec=10min
OnActiveSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
